# track-phone.conf
description "track phone websocket server"
start on runlevel [2345]
env DISPLAY=":0.0"
env HOME="/home/USERNAME"
start on (started lightdm)
exec /opt/track-phone/server.py
